# chef-contrail/cookbooks/redis/files/default/redis-query.conf
description "redis-query"

start on local-filesystems
stop on shutdown

setuid redis
setgid redis

exec /usr/bin/redis-server /etc/redis/redis-query.conf

respawn
