# chef-contrail/cookbooks/redis/files/default/redis-web.conf
description "redis-web"

start on local-filesystems
stop on shutdown

setuid redis
setgid redis

exec /usr/bin/redis-server /etc/redis/redis-web.conf

respawn
