# chef-contrail/cookbooks/nova/files/default/nova-save.conf
# Generated by chef
description "upstart script to save vm info before shutdown"

start on runlevel [06]

script
  [ -f /home/nova/tools/keystonecreds ] && . /home/nova/tools/keystonecreds
end script

task
exec /etc/init/nova-save.py
